# password-vault.service
[Unit]
Description=Password Vault - Secure Password Manager
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/project/experimental/password_saver/password_saver
Environment=PYTHONPATH=%h/project/experimental/password_saver/password_saver
Environment=FLASK_ENV=production
Environment=FLASK_SECRET=your_secure_secret_key_here_change_this
Environment=HTTPS=0
ExecStart=/usr/bin/python3 app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security settings (simplified for user service)
NoNewPrivileges=yes
PrivateTmp=yes

# Resource limits for optimization
MemoryLimit=64M
CPUQuota=50%
TasksMax=10

[Install]
WantedBy=default.target 